# roles/raspberry_pi/templates/etc/profile.d/sdr-pi.sh.j2
export PATH="{{ sdr_pi_folder }}/bin:${PATH}"