# roles/raspberry_pi/templates/etc/profile.d/sdr-pi.sh.j2
export PATH="${PATH}:{{ sdr_pi_folder }}/bin"