# roles/bettercap/templates/etc/profile.d/go-env.sh.j2
export PATH="${PATH}:{{ sdr_pi_folder }}/go/bin"
export GOPATH="{{ sdr_pi_folder }}/go"
