# roles/soapysdr_server/templates/etc/profile.d/ettus.sh.j2
export UHD_IMAGES_DIR=/usr/share/uhd/images